# salt/rhnet/custom-supervisor.conf
[program:rhnet]
command=npm start
directory=/home/rh-user/public_html
user=rh-user
